# ha-llm-addon/rootfs/etc/s6-overlay/s6-rc.d/llm-wrapper/finish
#!/usr/bin/env bash
# s6-overlay finish script — called when the run script exits or is killed.

set -euo pipefail

echo "[ha-llm] LLM Wrapper stopped."
