# ha-llm-addon/rootfs/etc/s6-overlay/s6-rc.d/llm-wrapper/finish
#!/usr/bin/env bash
# shellcheck source=/dev/null
source /usr/lib/bashio/bashio.sh

bashio::log.info "LLM Wrapper stopped."
