# ha-llm/rootfs/etc/s6-overlay/s6-rc.d/llm-wrapper/run
#!/usr/bin/with-contenv bashio
# vim: ft=bash
# shellcheck shell=bash

# ── Read options via bashio ───────────────────────────────────────────────────
OLLAMA_HOST=$(bashio::config 'ollama_host'    'http://homeassistant.local:11434')
MODEL=$(bashio::config        'model'          'llama3')
MAX_TOKENS=$(bashio::config   'max_tokens'     '1024')
CTX_WINDOW=$(bashio::config   'context_window' '8192')
HA_TOKEN=$(bashio::config     'ha_token'       '')
HA_URL=$(bashio::config       'ha_url'         'http://supervisor/core')
SYSTEM_PROMPT=$(bashio::config 'system_prompt' '')

# ── Export as HA_LLM_* for pydantic-settings ─────────────────────────────────
export HA_LLM_OLLAMA_HOST="$OLLAMA_HOST"
export HA_LLM_MODEL="$MODEL"
export HA_LLM_MAX_TOKENS="$MAX_TOKENS"
export HA_LLM_CONTEXT_WINDOW="$CTX_WINDOW"
export HA_LLM_HA_URL="$HA_URL"
export HA_LLM_SYSTEM_PROMPT="$SYSTEM_PROMPT"

# SUPERVISOR_TOKEN est injecté automatiquement par le Supervisor.
# Utilisé comme fallback si ha_token n'est pas renseigné dans les options.
if [[ -z "$HA_TOKEN" && -n "${SUPERVISOR_TOKEN:-}" ]]; then
    export HA_LLM_HA_TOKEN="$SUPERVISOR_TOKEN"
else
    export HA_LLM_HA_TOKEN="$HA_TOKEN"
fi

bashio::log.info "Starting LLM Wrapper — model=${MODEL}, ollama=${OLLAMA_HOST}, ha=${HA_URL}"

cd /app

exec python3 -m uvicorn app.main:app \
    --host 0.0.0.0 \
    --port 8080 \
    --log-level info
