# ha-ollama/rootfs/etc/s6-overlay/s6-rc.d/ollama/finish
#!/command/with-contenv bashio
# vim: ft=bash
# shellcheck shell=bash

bashio::log.info "Stopping Ollama server..."

pkill -TERM -x ollama 2>/dev/null || true
sleep 5

bashio::log.info "Ollama stopped."
