# ha-ollama/rootfs/etc/s6-overlay/s6-rc.d/ollama/run
#!/usr/bin/with-contenv bashio
# vim: ft=bash
# shellcheck shell=bash

# ── Read options via bashio ───────────────────────────────────────────────────
DATA_DIR=$(bashio::config  'data_dir'          '/data/ollama')
KEEP_ALIVE=$(bashio::config 'keep_alive'       '5m')
MAX_LOADED=$(bashio::config 'max_loaded_models' '1')
GPU_LAYERS=$(bashio::config 'gpu_layers'        '0')
NUM_THREADS=$(bashio::config 'num_threads'      '4')

# models est un tableau JSON — on le convertit en tableau bash
mapfile -t MODELS < <(bashio::config 'models' | jq -r '.[]' 2>/dev/null || true)

# ── Environment variables understood by Ollama ────────────────────────────────
export OLLAMA_MODELS="$DATA_DIR/models"
export OLLAMA_HOST="0.0.0.0:11434"
export OLLAMA_KEEP_ALIVE="$KEEP_ALIVE"
export OLLAMA_MAX_LOADED_MODELS="$MAX_LOADED"
export OLLAMA_NUM_PARALLEL="1"
export OLLAMA_NUM_THREAD="$NUM_THREADS"

if [[ "$GPU_LAYERS" -gt 0 ]]; then
    export OLLAMA_GPU_LAYERS="$GPU_LAYERS"
    bashio::log.info "GPU layers enabled: ${GPU_LAYERS}"
else
    export OLLAMA_NO_CUDA=1
    export OLLAMA_NO_GPU=1
fi

# ── Create data directory ─────────────────────────────────────────────────────
mkdir -p "$OLLAMA_MODELS"
bashio::log.info "Model storage: ${OLLAMA_MODELS}"

# ── Start Ollama server ───────────────────────────────────────────────────────
bashio::log.info "Starting Ollama server on ${OLLAMA_HOST} ..."
ollama serve &
OLLAMA_PID=$!

# ── Wait for the API to be ready (max 60 s) ───────────────────────────────────
bashio::log.info "Waiting for Ollama API..."
for i in $(seq 1 60); do
    if curl -sf http://localhost:11434/api/tags > /dev/null 2>&1; then
        bashio::log.info "Ollama API ready after ${i}s"
        break
    fi
    if [[ $i -eq 60 ]]; then
        bashio::log.error "Ollama API did not start within 60 seconds — aborting"
        kill "$OLLAMA_PID" 2>/dev/null || true
        exit 1
    fi
    sleep 1
done

# ── Pull configured models if not already present ─────────────────────────────
if [[ ${#MODELS[@]} -eq 0 ]]; then
    bashio::log.info "No models configured — skipping pull"
else
    for model in "${MODELS[@]}"; do
        bashio::log.info "Checking model: ${model}"
        if ollama list 2>/dev/null | grep -q "^${model}"; then
            bashio::log.info "  already present"
        else
            bashio::log.info "  pulling ${model} (this may take several minutes)..."
            if ollama pull "$model"; then
                bashio::log.info "  ${model} pulled successfully"
            else
                bashio::log.warning "  Failed to pull ${model} — continuing"
            fi
        fi
    done
fi

bashio::log.info "Ollama ready. Loaded: $(ollama list 2>/dev/null | tail -n +2 | awk '{print $1}' | tr '\n' ' ')"

wait "$OLLAMA_PID"
